# ops/seo-audit-runner/deploy/systemd/seo-audit-runner.service
# seo-audit-runner.service — one full audit run (`run --all`).
# Started by seo-audit-runner.timer. SHIPS DISABLED — enabling is an
# explicit administrator action gated by docs/PRODUCTION_GATES.md.
[Unit]
Description=SEO audit runner — full audit of all projects
Documentation=file:/opt/seo-audit-runner/current/README.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=seo-runner
Group=seo-runner
EnvironmentFile=/etc/seo-audit-runner/runner.env
ExecStart=/usr/local/bin/seo-audit-runner run --all
Nice=10

# Timeouts and shutdown: audits poll long-running work; SIGTERM triggers the
# runner's graceful abort (lock released, run journaled as ABORTED).
TimeoutStartSec=2h
KillSignal=SIGTERM
TimeoutStopSec=90

# Resource limits
MemoryHigh=1G
MemoryMax=2G
TasksMax=128

# Hardening: the runner may write ONLY its own state and log directories.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/seo-audit-runner /var/log/seo-audit-runner
RuntimeDirectory=seo-audit-runner
RuntimeDirectoryMode=0750
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
UMask=0027
