# ops/seo-audit-runner/deploy/systemd/seo-audit-runner.timer
# seo-audit-runner.timer — daily full audit at 03:00 Africa/Cairo.
# SHIPS DISABLED. Never enable before the production gates pass
# (docs/PRODUCTION_GATES.md Gates 5-7). Choose EITHER this classic daily
# timer OR runner-managed schedules via seo-runner-tick.timer for full
# audits — never both (duplicate audits and alert noise).
[Unit]
Description=Daily SEO audit (03:00 Africa/Cairo)

[Timer]
# Timezone pinned in the calendar expression (systemd >= 235). Alternative:
# set the host timezone to Africa/Cairo and use `*-*-* 03:00:00`.
OnCalendar=*-*-* 03:00:00 Africa/Cairo
RandomizedDelaySec=300
Persistent=true
Unit=seo-audit-runner.service

[Install]
WantedBy=timers.target
