# ops/seo-audit-runner/deploy/systemd/seo-runner-retry.timer
# seo-runner-retry.timer — hourly Slack notification retry.
# SHIPS DISABLED. Enable together with (or after) the audit scheduling
# decision — see docs/PRODUCTION_GATES.md.
[Unit]
Description=Hourly SEO runner notification retry

[Timer]
OnCalendar=hourly
RandomizedDelaySec=120
Persistent=true
Unit=seo-runner-retry.service

[Install]
WantedBy=timers.target
